# libssl.service
[Unit]
Description=LibSSL_GSurf
[Service]
ExecStart=/gsurf/serverSSL -s sitef
[Install]
WantedBy=multi-user.target